# infra/backup/slate-backup.timer
[Unit]
Description=Run the Slate DB backup daily

[Timer]
# Daily at 03:30 UTC. Persistent=true runs a missed backup after downtime.
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
Unit=slate-backup.service

[Install]
WantedBy=timers.target
